# systemd/media_processor.service
[Unit]
Description=Media Processing Service (Download -> Process -> Upload)
After=network-online.target

[Service]
Type=oneshot
# The working directory is the app's root, so the script can find its config.
WorkingDirectory=/root/mediaflow
# The user must have permissions to execute this script.
ExecStart=/root/mediaflow/run_all.sh

[Install]
WantedBy=multi-user.target
