# systemd/media_processor.service
[Unit]
Description=Media Processing Service (Download -> Process -> Upload)
After=network-online.target

[Service]
Type=oneshot
# The WorkingDirectory should be the root of your cloned repository,
# which the deploy_ui.sh script places at /root/mediaflow
WorkingDirectory=/root/mediaflow
# The ExecStart should be the full path to your script inside that directory
ExecStart=/root/mediaflow/run_all.sh
