# systemd/media_processor.service
[Unit]
Description=Media Processing Service (Download -> Process -> Upload)
After=network-online.target

[Service]
Type=oneshot
# IMPORTANT: Update the path to your run_all.sh script
ExecStart=/usr/local/bin/run_all.sh
# IMPORTANT: Update the user to the one who owns the project files
User=your_username
